# root/etc/services.d/papermerge-importer/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/papermerge-importer || exit

exec \
	s6-setuidgid abc ./papermerge-importer.sh
